# ha-addon/rootfs/etc/services.d/bot/run
#!/command/with-contenv bashio
set -euo pipefail

cd /opt/bot

if [ -f /data/bot/.env ]; then
    set -a
    source /data/bot/.env
    set +a
fi

if [ -f /var/run/ha_bot_enable_gateway ] && grep -qi "true" /var/run/ha_bot_enable_gateway; then
    PORT="$(cat /var/run/ha_bot_gateway_port 2>/dev/null || echo 6443)"
    bashio::log.info "Starting Home Assistant gateway on port ${PORT}"
    python3 -m uvicorn tools.ha_gateway:app --host 0.0.0.0 --port "${PORT}" --workers 1 --log-level info &
fi

if [ -f /var/run/ha_bot_enable_log_gateway ] && grep -qi "true" /var/run/ha_bot_enable_log_gateway; then
    TRADE_PORT="$(cat /var/run/ha_bot_trade_port 2>/dev/null || echo 6442)"
    ERROR_PORT="$(cat /var/run/ha_bot_error_port 2>/dev/null || echo 6441)"
    bashio::log.info "Starting trade log gateway on port ${TRADE_PORT}"
    python3 -m uvicorn tools.log_gateway:trade_app --host 0.0.0.0 --port "${TRADE_PORT}" --workers 1 --log-level info &
    bashio::log.info "Starting error log gateway on port ${ERROR_PORT}"
    python3 -m uvicorn tools.log_gateway:error_app --host 0.0.0.0 --port "${ERROR_PORT}" --workers 1 --log-level info &
fi

bashio::log.info "Starting trading bot"
exec python3 -m bot
