# ha-addon/rootfs/etc/services.d/bot/run
#!/command/with-contenv bashio
set -euo pipefail

cd /opt/bot

if [ -f /data/bot/.env ]; then
    set -a
    source /data/bot/.env
    set +a
fi

if [ -f /var/run/ha_bot_enable_gateway ] && grep -qi "true" /var/run/ha_bot_enable_gateway; then
    PORT="$(cat /var/run/ha_bot_gateway_port 2>/dev/null || echo 6443)"
    bashio::log.info "Starting Home Assistant gateway on port ${PORT}"
    python3 -m uvicorn tools.ha_gateway:app --host 0.0.0.0 --port "${PORT}" --workers 1 --log-level info &
fi

bashio::log.info "Starting trading bot"
exec python3 -m bot
