# home_assistant/bithumb_kis_bot/rootfs/etc/services.d/bot/run
#!/command/with-contenv bashio
set -euo pipefail

cd /opt/bot

if [ -f /data/bot/.env ]; then
    set -a
    source /data/bot/.env
    set +a
fi

if [ -f /var/run/ha_bot_enable_gateway ] && grep -q "true" /var/run/ha_bot_enable_gateway; then
    bashio::log.info "Starting Home Assistant gateway on port 8080"
    python -m tools.ha_gateway --host 0.0.0.0 --port 8080 &
fi

bashio::log.info "Starting trading bot"
exec python -m bot.bithumb_bot
